# deploy/systemd/kinetix-web.service
[Unit]
Description=Kinetix Web (Next.js)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/clickup/kinetix/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
# Do NOT use ExecStartPre with pkill -f "next dev" — the shell cmdline contains
# that string and pkill kills itself, causing "fatal signal to control process".
ExecStart=/opt/clickup/kinetix/frontend/node_modules/.bin/next start -p 3000 -H 127.0.0.1
Restart=always
RestartSec=5
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
